# klereo/rootfs/etc/services.d/klereo/run
#!/usr/bin/bash

echo "$(date): Klereo add-on starting up..."
echo "$(date): Current user: $(whoami)"
echo "$(date): Current directory: $(pwd)"
echo "$(date): Environment variables:"
env | grep -E "(SUPERVISOR|HOME|PATH)" | head -10

# Load configuration from Home Assistant
CONFIG_PATH="/data/options.json"

echo "$(date): Looking for configuration at: $CONFIG_PATH"
if [ ! -f "$CONFIG_PATH" ]; then
    echo "$(date): ERROR: Configuration file not found at $CONFIG_PATH"
    echo "$(date): Available files in /data:"
    ls -la /data/ 2>/dev/null || echo "$(date): /data directory not accessible"
    exit 1
fi

echo "$(date): Configuration file found, parsing..."

# Parse configuration using jq
klereo_username=$(jq -r '.klereo_username // empty' "$CONFIG_PATH")
klereo_password=$(jq -r '.klereo_password // empty' "$CONFIG_PATH")
update_interval=$(jq -r '.update_interval // 600' "$CONFIG_PATH")
log_level=$(jq -r '.log_level // "info"' "$CONFIG_PATH")

# Export configuration as environment variables
export KLEREO_USERNAME="${klereo_username}"
export KLEREO_PASSWORD="${klereo_password}"
export UPDATE_INTERVAL="${update_interval}"
export LOG_LEVEL="${log_level^^}"

# Home Assistant add-ons have automatic access to the supervisor API
export HOMEASSISTANT_URL="http://supervisor/core"
# The supervisor token is automatically available
export HOMEASSISTANT_TOKEN="${SUPERVISOR_TOKEN}"

# Log startup information
echo "$(date): Starting Klereo Pool Manager Add-on..."
echo "$(date): Log level: ${log_level}"
echo "$(date): Update interval: ${update_interval}s"
echo "$(date): Home Assistant integration: automatic via supervisor"

# Validate configuration
if [ -z "${klereo_username}" ]; then
    echo "$(date): ERROR: Klereo username is required!"
    exit 1
fi

if [ -z "${klereo_password}" ]; then
    echo "$(date): ERROR: Klereo password is required!"
    exit 1
fi

# Change to working directory
cd /usr/bin || exit 1

# Start the application
echo "$(date): Starting Klereo Pool Manager service..."
exec python3 /usr/bin/klereo